# systemd/meshymcmapface-agent.service
[Unit]
Description=MeshyMcMapface Agent - Meshtastic Network Monitor
Documentation=https://github.com/heymex/MeshyMcMapface
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=meshyuser
Group=meshyuser
SupplementaryGroups=dialout
WorkingDirectory=/opt/meshymcmapface
PIDFile=/run/meshymcmapface/agent.pid
ExecStart=/opt/meshymcmapface/venv/bin/python3 /opt/meshymcmapface/mmm-agent-modular.py --config /etc/meshymcmapface/agent.ini --log-level INFO --log-file /var/log/meshymcmapface/agent.log --pid-file /run/meshymcmapface/agent.pid
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
RuntimeDirectory=meshymcmapface
RuntimeDirectoryMode=0755

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/meshymcmapface /var/lib/meshymcmapface /etc/meshymcmapface
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Allow access to serial devices
DeviceAllow=/dev/ttyUSB* rw
DeviceAllow=/dev/ttyACM* rw
DeviceAllow=char-usb_device rw
DevicePolicy=closed

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Environment
Environment=PYTHONPATH=/opt/meshymcmapface/src
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshymcmapface-agent

[Install]
WantedBy=multi-user.target